# deploy/systemd/OpenTron.service
[Unit]
Description=OpenTron API Server
After=network.target

[Service]
Type=simple
User=opentron
WorkingDirectory=/opt/opentron
ExecStart=/opt/opentron/.venv/bin/tron serve --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
Environment=HOME=/opt/opentron
# Binding 0.0.0.0 requires authentication. This file MUST exist and contain:
#   OPENTRON_API_KEY=<key>      (generate one: `tron auth generate-key`)
# It is not prefixed with "-", so the unit fails to start if the file is
# missing — preventing an accidentally unauthenticated public server.
EnvironmentFile=/etc/opentron/env

[Install]
WantedBy=multi-user.target
